# templates/etc/systemd/system/statsd.service
[Unit]
Description=Statsd Service
After=syslog.target network.target

[Service]
Type=simple
Group=statsd
User=statsd
ExecStart=node test

[Install]
WantedBy=multi-user.target
